# deploy/arke-sync.service
[Unit]
Description=Arke Sync — rclone -> digest/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/arke/arke.env
ExecStart=arke-sync
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
