# deploy/arke-sync.service
[Unit]
Description=Arke Sync — rclone -> digest/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arke
Group=arke
WorkingDirectory=/opt/arke/src
EnvironmentFile=/opt/arke/config/arke.env
ExecStart=/opt/arke/src/.venv/bin/arke-sync
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
