# deploy/arke-mail.service
[Unit]
Description=Arke Mail — M365 email adapter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arke
Group=arke
WorkingDirectory=/opt/arke/src
EnvironmentFile=/opt/arke/config/arke.env
ExecStart=/opt/arke/src/.venv/bin/arke-mail
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60

# arke-mail may spawn `cloudflared` as a child (CLOUDFLARED_TUNNEL=<name>).
# control-group kill is the default — it reaps the tunnel along with the parent.

[Install]
WantedBy=multi-user.target
